# config/autostart/display.desktop
[Desktop Entry]
Name=Display
Comment=Tweak display settings
Exec=/home/ryan/scripts/linux/tweak-display.sh
Terminal=false
Type=Application
